# udev/99-yubikey.rules
ACTION=="remove", ENV{ID_VENDOR_ID}=="1050", ENV{ID_MODEL_ID}=="0407", ENV{ID_SERIAL_SHORT}=="000xxxxxxx", RUN+="/srv/yubikey-lock"
